# user.overlay/etc/init.d/bootstrap
#!/sbin/openrc-run


depend() {
        need networking
}

start() {
        ebegin "Starting SD card bootstrapper"
        /bootstrap &
        eend $?
}

